# deploy-service/ping-rest-api.service
[Unit]
Description=Evaporator REST API Service
After=network.target

[Service]
ExecStart=/home/eva-admin/evaporator-rest-api/app.exe
WorkingDirectory=/home/eva-admin/evaporator-rest-api
User=eva-admin
Restart=always
RestartSec=10
Environment=PORT=8016

[Install]
WantedBy=multi-user.target
